# crontab
0 3 * * * /home/dang/.scripts/dsync-cron
